# runtime/service/npud/org.tizen.npud.service.in
[Unit]
Description=Launch npu daemon

[Service]
Type=dbus
BusName=org.tizen.npud
User=owner
Group=users
Capabilities=cap_dac_override=i
SecureBits=keep-caps
Restart=always
RestartSec=3
ExecStart=@BINDIR@/npud
